# salt-minion-centos/files/etc/supervisord.d/docker.conf
[program:docker]
command=dockerd
stdout_logfile=/var/log/docker-stdout.log
autostart=true
startsecs=30