# salt-minion-centos/files/etc/supervisord.d/salt-minion.conf
[program:salt-minion]
command=salt-minion
stdout_logfile=/var/log/salt/salt-minion-stdout.log
autostart=true
autorestart=true
startsecs=0
