# salt-master/files/etc/supervisord.d/salt-master.conf
[program:salt-master]
command=salt-master
stdout_logfile=/var/log/salt/salt-master-stdout.log
autostart=true