# salt-master/files/etc/supervisord.d/salt-api.conf
[program:salt-api]
command=salt-api
stdout_logfile=/var/log/salt/salt-api-stdout.log
autostart=true
